# etc/cron.d/nfs_automount
* * * * *	root	/usr/local/bin/nfs_automount
